# config/includes.chroot/etc/systemd/system/nas-pro.service
[Unit]
Description=NAS-PRO Web Interface and API
After=network.target nginx.service
Wants=network.target

[Service]
Type=simple
User=naspro
Group=naspro
WorkingDirectory=/var/www/html
ExecStart=/usr/bin/node /var/www/html/server/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
