# config/systemd/nas-pro.service
[Unit]
Description=NAS-PRO Desktop

[Service]
ExecStart=/usr/local/bin/nas-os
Restart=always

[Install]
WantedBy=multi-user.target
